# tod_wallpaper.service
[Unit]
Description=Time Of Day Wallpaper

[Service]
ExecStart=/usr/bin/tod_wallpaper

[Install]
WantedBy=multi-user.target
